# rootfs/etc/s6-overlay/s6-rc.d/app/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -e

cd /app

exec node build/index.js